# launchd/com.finance-agent.analysis.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.finance-agent.analysis</string>

    <key>Comment</key>
    <string>Finance Agent — Path A: Scheduled portfolio analysis (weekdays 9:00 &amp; 13:00)</string>

    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/bin/python3</string>
      <string>/Users/zeru/.coding-agent/worktrees/finance-agent/issue-10-launchd/scripts/run-analysis.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/zeru/.coding-agent/worktrees/finance-agent/issue-10-launchd</string>

    <key>StartCalendarInterval</key>
    <array>
      <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
      <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>13</integer>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
      <dict>
        <key>Weekday</key>
        <integer>2</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
      <dict>
        <key>Weekday</key>
        <integer>2</integer>
        <key>Hour</key>
        <integer>13</integer>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
      <dict>
        <key>Weekday</key>
        <integer>3</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
      <dict>
        <key>Weekday</key>
        <integer>3</integer>
        <key>Hour</key>
        <integer>13</integer>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
      <dict>
        <key>Weekday</key>
        <integer>4</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
      <dict>
        <key>Weekday</key>
        <integer>4</integer>
        <key>Hour</key>
        <integer>13</integer>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
      <dict>
        <key>Weekday</key>
        <integer>5</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
      <dict>
        <key>Weekday</key>
        <integer>5</integer>
        <key>Hour</key>
        <integer>13</integer>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
    </array>

    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/zeru/.coding-agent/worktrees/finance-agent/issue-10-launchd/logs/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/zeru/.coding-agent/worktrees/finance-agent/issue-10-launchd/logs/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>/Users/zeru</string>
      <key>ALLOWED_OPEN_ID</key>
      <string>ou_c51af1689cf9f949f03d7ca822029c9e</string>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/Users/zeru/.bun/bin:/Users/zeru/.deskclaw/node/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
  </dict>
</plist>
